# docker/deploy/etc/s6-overlay/s6-rc.d/websockets/run
#!/command/with-contenv bash
s6-setuidgid webuser php $WEBUSER_HOME/artisan websockets:serve --no-ansi -q